# user/ccpd.timer

[Unit]
Description=Check crypto prices

[Timer]
OnBootSec=1min
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
