# user/ccpd.service

[Unit]
Description=ip checker

[Service]
ExecStart=/home/lef/work/ccpd

[Install]
WantedBy=multi-user.target
