# user/ngrok.service

[Unit]
Description=ngrok service

[Service]
ExecStart=ngrok tcp 22

[Install]
WantedBy=default.target
