# user/random_background.service

[Unit]
Description=Change desktop background every hour

[Service]
ExecStart=/home/lef/.local/bin/setbg

[Install]
WantedBy=multi-user.target
